# supervisor/aawt.conf
[program:aawt]
command = python aawt.py
directory = /opt/aawt
autostart = true
autorestart = true
