# supervisor/atw.conf
[program:atw]
command = python atw.py
directory = /opt/atw
autostart = true
autorestart = true
